# usr/lib/systemd/system/ssh.service.d/99_sysmaint.conf
[Unit]
After=sysmaint-boot.target

[Install]
WantedBy=sysmaint-boot.target
