# usr/lib/systemd/system/usbguard-dbus.service.d/99_sysmaint.conf
[Install]
WantedBy=sysmaint-boot.target
